# srv/fai/config/files/archives/CIP/etc/xdg/autostart/guake.desktop
[Desktop Entry]
Name=Autostart Guake
Exec=/usr/bin/guake &
Terminal=false
Type=Application
