# srv/fai/config/files/archives/CIP/etc/xdg/autostart/script-launcher.desktop
[Desktop Entry]
Name=Launch CIP scripts
Exec=/usr/local/bin/script-launcher &
Terminal=false
Type=Application
